# demo/assets/supervisor.d/supervisord.conf
[supervisord]
user=root
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:sbdaemon]
command=/opt/sb/sb -d
autostart=true
autorestart=true
stdout_logfile=/var/log/sb.log
stderr_logfile=/var/log/sb.log